# root/etc/s6-overlay/s6-rc.d/nexus/run
#!/command/with-contenv sh

cd "${NEXUS_HOME}"
exec s6-setuidgid appuser bin/nexus run
