# root/etc/s6-overlay/s6-rc.d/nexus-limits/run
#!/command/with-contenv sh

echo "Adjusting system limits"
ulimit -n 65536 -u 4096

