# root/etc/s6-overlay/s6-rc.d/nexus-javaprefs/run
#!/command/with-contenv sh

  if [ ! -f ${NEXUS_DATA}/javaprefs ]; then
    echo "java user prefs file is missing, creating empty one"
    touch ${NEXUS_DATA}/javaprefs
  fi
